# zhaws/rootfs/etc/services.d/zhaws/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start ZHAWS service for Zigbee radio
# ==============================================================================

# Send out discovery information to Home Assistant
./discovery &

# Run daemon, passing external config directory in as environment variable
python3 -m zhaws.server --config /etc/zhaws_config.json
